# freechat.service
[Unit]
Description=FreeChat - AI chat CLI tool
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/freechat
ExecStart=/usr/bin/python3 /path/to/freechat/freechat.py
Restart=on-failure
RestartSec=5s

# Environment variables (uncomment and set your API keys)
# Environment=OPENAI_API_KEY=your_openai_api_key
# Environment=OPENROUTER_API_KEY=your_openrouter_api_key
# Environment=GEMINI_API_KEY=your_gemini_api_key
# Environment=ANTHROPIC_API_KEY=your_anthropic_api_key
# Environment=MISTRAL_API_KEY=your_mistral_api_key

[Install]
WantedBy=multi-user.target
